# recipes-support/linht-utils/files/linht-volume-ctrl.service
[Unit]
Description=LinHT Volume Control Service

[Service]
Type=simple
ExecStart=/usr/bin/volume_ctrl.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target